# ROOT/etc/profile.d/profile_dotfiles_init_tenv.sh
#!/bin/bash
# https://github.com/tenv/tenv
export TENV_AUTO_INSTALL=true
