# ROOT/etc/profile.d/profile_dotfiles_init_go.sh
#!/bin/bash
echo 'export GOPATH=$HOME' > $HOME/.go
echo 'export PATH=$PATH:$HOME/go/bin' >> $HOME/.go
echo 'export PATH=$PATH:$GOPATH/bin' >> $HOME/.go

source $HOME/.go
