# ROOT/etc/profile.d/profile_dotfiles_init_aws.sh
#!/bin/bash
# https://docs.aws.amazon.com/cli/latest/userguide/cli-configure-completion.html
if command -v aws_completer >/dev/null 2>&1; then
    complete -C "$HOME/.asdf/shims/aws_completer" aws
fi

fi
