# ROOT/etc/profile.d/profile_dotfiles_init_direnv.sh
#!/bin/bash

# disable direnv hook
#export EDITOR=vim
#eval $(direnv hook bash)
