# ROOT/etc/profile.d/profile_dotfiles_aqua.sh
#!/bin/bash
if [[ -d "${HOME}/.local/share/aquaproj-aqua/bin" ]]; then
    export PATH="${HOME}/.local/share/aquaproj-aqua/bin:$PATH"
fi
