# root/etc/profile.d/profile_ansible_init_kubectl.sh
if kubectl version ; then
    source <(kubectl completion bash)
fi
