# ROOT/etc/profile.d/profile_dotfiles_init_homebrew.sh
#!/bin/bash
# https://docs.brew.sh/Homebrew-on-Linux
test -d "${HOME}/.linuxbrew" && eval "$(${HOME}/.linuxbrew/bin/brew shellenv)"
